# scripts/com.gyturas.csi1000.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.gyturas.csi1000.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Library/Frameworks/Python.framework/Versions/3.11/bin/python3</string>
    <string>-m</string>
    <string>csi1000.engine.daily_dashboard</string>
  </array>

  <key>WorkingDirectory</key><string>/Users/yilinzhou/Desktop/中证1000择时</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key><string>/Users/yilinzhou/Desktop/中证1000择时/src</string>
    <key>PATH</key><string>/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin</string>
    <key>MPLBACKEND</key><string>Agg</string>
    <key>CSI1000_NO_OPEN</key><string>1</string>
  </dict>

  <!-- 周一至周五 15:30(A股15:00收盘,留30分钟等数据落地) -->
  <key>StartCalendarInterval</key>
  <array>
    <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
    <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
  </array>

  <key>StandardOutPath</key><string>/Users/yilinzhou/Desktop/中证1000择时/logs/daily.log</string>
  <key>StandardErrorPath</key><string>/Users/yilinzhou/Desktop/中证1000择时/logs/daily.err.log</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
